# packages/exedev/s6-overlay/rootfs/etc/s6-overlay/s6-rc.d/exedev-authorized-keys/up
/command/with-contenv sh /etc/s6-overlay/scripts/exedev-authorized-keys
