# packages/exedev/s6-overlay/rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv sh
# sshd refuses to start unless invoked with an absolute path; resolve it on PATH.
exec "$(command -v sshd)" -D -e -f /etc/ssh/sshd_config
