# packages/exedev/s6-overlay/rootfs/etc/s6-overlay/s6-rc.d/nix-daemon/run
#!/command/with-contenv sh
exec nix-daemon --daemon
