# packages/exedev/s6-overlay/rootfs/etc/s6-overlay/s6-rc.d/sshd-keygen/up
/command/with-contenv sh /etc/s6-overlay/scripts/sshd-keygen
